# debian/DEBIAN/postinst
#!/bin/bash
set -e

# Reload systemd to recognize new services (skip if systemctl not available)
if command -v systemctl >/dev/null 2>&1 && systemctl >/dev/null 2>&1; then
    systemctl daemon-reload || true
    # Disable service by default (it should not auto-start)
    systemctl disable luoking-box.service >/dev/null 2>&1 || true
fi

# Create log directory
mkdir -p /var/log/luoking-box
chmod 755 /var/log/luoking-box

# Create sing-box-config directory if it doesn't exist
mkdir -p /etc/luoking-box/sing-box-config
chmod 755 /etc/luoking-box/sing-box-config

# Set permissions for config files (if they exist from package)
if [ -f /etc/luoking-box/config.json ]; then
    chmod 600 /etc/luoking-box/config.json
fi

if [ -f /etc/luoking-box/sing-box-config/default.json ]; then
    chmod 600 /etc/luoking-box/sing-box-config/default.json
fi

# Create default config.json if it doesn't exist (for upgrades)
if [ ! -f /etc/luoking-box/config.json ]; then
    echo '{"active_config": "default"}' > /etc/luoking-box/config.json
    chmod 600 /etc/luoking-box/config.json
    echo "Created default config.json at /etc/luoking-box/config.json"
fi

# Create default configuration file if it doesn't exist (for upgrades)
if [ ! -f /etc/luoking-box/sing-box-config/default.json ]; then
    cat > /etc/luoking-box/sing-box-config/default.json << 'EOF'
{
    "log": {
        "disabled": true,
        "level": "info",
        "timestamp": true
    },
    "dns": {
        "servers": [],
        "rules": [],
        "final": "",
        "strategy": "",
        "disable_cache": false,
        "disable_expire": false,
        "independent_cache": false,
        "cache_capacity": 0,
        "reverse_mapping": false
    },
    "inbounds": [],
    "outbounds": [],
    "route": {
        "rules": [],
        "rule_set": [],
        "default_domain_resolver": "",
        "final": "",
        "auto_detect_interface": false
    }
}
EOF
    chmod 600 /etc/luoking-box/sing-box-config/default.json
    echo "Created default configuration at /etc/luoking-box/sing-box-config/default.json"
    echo "Please configure it before starting the service."
fi

# Auto-configure shell integration for users
add_shell_integration() {
    local user_home="$1"
    local shell_rc="$2"
    
    if [ -z "$shell_rc" ] || [ ! -w "$(dirname "$shell_rc")" ]; then
        return 1  # Cannot write to directory
    fi
    
    # Create file if it doesn't exist
    if [ ! -f "$shell_rc" ]; then
        touch "$shell_rc"
        chmod 644 "$shell_rc"
        # Try to set ownership if we know the user
        if [ -n "$INSTALL_USER" ] && [ "$INSTALL_USER" != "root" ]; then
            chown "$INSTALL_USER:$INSTALL_USER" "$shell_rc" 2>/dev/null || true
        fi
    fi
    
    local marker="# luoking-box shell integration"
    if grep -qF "$marker" "$shell_rc" 2>/dev/null; then
        return 0  # Already added
    fi
    
    echo "" >> "$shell_rc"
    echo "$marker" >> "$shell_rc"
    echo "[ -f /etc/profile.d/luoking-box.sh ] && source /etc/profile.d/luoking-box.sh" >> "$shell_rc"
    
    return 0
}

# Detect and configure shell integration for current user
if [ -n "$SUDO_USER" ]; then
    # Running with sudo, use SUDO_USER
    INSTALL_USER="$SUDO_USER"
elif [ "$(id -u)" -ne 0 ] && [ -n "$USER" ]; then
    # Running as regular user
    INSTALL_USER="$USER"
else
    # Running as root without sudo, try to detect from environment
    INSTALL_USER="${SUDO_USER:-${USER:-}}"
fi

SHELL_INTEGRATION_ADDED=false

if [ -n "$INSTALL_USER" ] && [ "$INSTALL_USER" != "root" ]; then
    USER_HOME=$(eval echo ~"$INSTALL_USER" 2>/dev/null || getent passwd "$INSTALL_USER" 2>/dev/null | cut -d: -f6 || echo "")
    
    if [ -n "$USER_HOME" ] && [ -d "$USER_HOME" ]; then
        # Detect shell and add to appropriate config file
        USER_SHELL=$(getent passwd "$INSTALL_USER" 2>/dev/null | cut -d: -f7 || echo "")
        
        if [[ "$USER_SHELL" == *"zsh"* ]]; then
            if add_shell_integration "$USER_HOME" "$USER_HOME/.zshrc"; then
                SHELL_INTEGRATION_ADDED=true
            fi
        elif [[ "$USER_SHELL" == *"bash"* ]] || [ -z "$USER_SHELL" ]; then
            # Default to bash if shell is unknown
            if add_shell_integration "$USER_HOME" "$USER_HOME/.bashrc"; then
                SHELL_INTEGRATION_ADDED=true
            fi
        fi
        
        # Also add to .profile for login shells (bash) - but only if .bashrc doesn't exist
        if [ ! -f "$USER_HOME/.bashrc" ] && [ -f "$USER_HOME/.profile" ]; then
            add_shell_integration "$USER_HOME" "$USER_HOME/.profile"
        fi
    fi
fi

# For system-wide installation, /etc/profile.d/luoking-box.sh will be automatically
# loaded for login shells. For non-login interactive shells, users need to have
# it in their ~/.bashrc or ~/.zshrc (which we just added above).

echo ""
echo "luoking-box has been installed successfully!"
echo ""
echo "To start the service (for local proxy):"
echo "  sudo systemctl start luoking-box"
echo ""
echo "Configuration files:"
echo "  Main config: /etc/luoking-box/config.json"
echo "  Config directory: /etc/luoking-box/sing-box-config/"
echo ""
echo "To switch configuration, edit /etc/luoking-box/config.json and change 'active_config'"
echo "then restart the service: sudo systemctl restart luoking-box"
echo ""
if [ "$SHELL_INTEGRATION_ADDED" = true ] && [ -n "$INSTALL_USER" ]; then
    echo "✓ Shell integration has been automatically configured for user: $INSTALL_USER"
    echo "  You can now directly use:"
    echo "    luoking-box enable session docker"
    echo "    luoking-box clear session docker"
    echo ""
    echo "  Note: Open a new terminal or run 'source ~/.bashrc' (or 'source ~/.zshrc') to activate."
    echo "  Or simply run 'luoking-box enable session' - it will auto-load shell integration."
else
    echo "Shell integration:"
    echo "  Shell integration will be auto-loaded when you run 'luoking-box enable session'."
    echo "  Or add this to your ~/.bashrc or ~/.zshrc for permanent activation:"
    echo "    [ -f /etc/profile.d/luoking-box.sh ] && source /etc/profile.d/luoking-box.sh"
fi

