# debian/DEBIAN/prerm
#!/bin/bash
set -e

# Stop service before removal (skip if systemctl not available)
if command -v systemctl >/dev/null 2>&1 && systemctl >/dev/null 2>&1; then
    if systemctl is-active --quiet luoking-box.service 2>/dev/null; then
        systemctl stop luoking-box.service || true
    fi
    # Disable service
    systemctl disable luoking-box.service >/dev/null 2>&1 || true
fi

